# config-overrides/includes.chroot/etc/skel/.config/systemd/user/swaybg.service
[Unit]
Description=Idle manager for Wayland
Documentation=man:swayidle(1)
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'swaymsg output \* bg $(gsettings get org.gnome.desktop.background picture-uri | cut -c 9- | rev | cut -c 2- | rev) fill'

[Install]
WantedBy=sway-session.target
