# config/includes.chroot/etc/skel/.config/autostart/xfce4-notes-autostart.desktop
[Desktop Entry]
Hidden=true
Exec=xfce4-notes
Terminal=false
