# config-overrides/includes.chroot/etc/rawfs/etc/skel/.config/systemd/user/swayidle.service
[Unit]
Description=Idle manager for Wayland
Documentation=man:swayidle(1)
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/swayidle -w before-sleep swaylock

[Install]
WantedBy=sway-session.target
